# docker/ue4_game_exec/etc/supervisor/conf.d/noVNC.conf
[program:noVNC]
command=bash -l -c "./utils/launch.sh"
environment=HOME=/home/ue4,USER=ue4
directory=/opt/noVNC
autorestart=true
stdout_logfile=/var/log/supervisor/noVNC.log
redirect_stderr=true
